# deploy/delta_api.service
# ==========================================
#  Acrome Delta API — Systemd Service
# ==========================================
#
#  This file is OPTIONAL. Use it only if you want the API to start
#  automatically every time the Raspberry Pi boots up.
#
#  Steps:
#  1. Update the paths below to match your system.
#  2. Copy this file:
#       sudo cp delta_api.service /etc/systemd/system/
#  3. Enable and start the service:
#       sudo systemctl daemon-reload
#       sudo systemctl enable delta_api.service
#       sudo systemctl start delta_api.service
#
#  Note: For most use cases, 'delta start' / 'delta stop' commands
#        (installed via setup.sh) are sufficient. This service file
#        is for advanced users who need auto-start on boot.
# ==========================================

[Unit]
Description=Acrome Delta Robot API Service
After=network.target

[Service]
Type=simple
# Set your Raspberry Pi username (usually 'pi')
User=pi
# Set the full path to the API directory:
WorkingDirectory=/home/pi/Delta-App/api
ExecStart=/home/pi/Delta-App/api/venv/bin/python3 /home/pi/Delta-App/api/app.py

# Restart on failure after 5 seconds
Restart=on-failure
RestartSec=5

# Flask environment variables
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
